# php-fpm/orocommerce/context/etc/supervisord.d/oro-web-socket.ini
[program:oro_web_socket]
command=php ./bin/console gos:websocket:server --env=prod
numprocs=1
autostart=true
autorestart=true
directory=/var/www/html
user=www-data
redirect_stderr=true