# php-fpm/orocommerce/context/etc/supervisord.d/oro-message-consumer.ini
[program:oro_message_consumer]
command=php ./bin/console oro:message-queue:consume --env=prod
process_name=%(program_name)s_%(process_num)02d
numprocs=5
autostart=true
autorestart=true
directory=/var/www/html
user=www-data
redirect_stderr=true